# systemd/thermal-monitor.service
# thermal-monitor.service — single collection cycle, triggered by the timer.
#
# Installation:
#   1. Place the project at /opt/thermal_monitor (or override paths below).
#      No service user is needed — DynamicUser=yes allocates a transient
#      UID/GID per invocation.
#   2. Edit /etc/thermal_monitor/thermal_monitor.yaml (copy from the example).
#      Point logging.db_file into /var/lib/thermal_monitor/ so the SQLite
#      reading log also lands in StateDirectory:
#          logging:
#            db_file: /var/lib/thermal_monitor/readings.db
#      (The alert-cooldown state_file is already handled for you via the
#      $THERMAL_MONITOR_STATE_FILE env var set below — PrivateTmp=yes would
#      otherwise isolate /tmp per run and reset the cooldown every cycle.)
#   3. Install the units:
#          sudo install -m 644 systemd/thermal-monitor.service /etc/systemd/system/
#          sudo install -m 644 systemd/thermal-monitor.timer   /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now thermal-monitor.timer
#   4. Watch it run:
#          journalctl -u thermal-monitor.service -f
#          journalctl -u thermal-monitor.service -p warning    # only WARN/CRIT
#
# Tune paths or override settings without editing this file by running:
#      sudo systemctl edit thermal-monitor.service
# which creates a drop-in at /etc/systemd/system/thermal-monitor.service.d/.

[Unit]
Description=Equipment-room thermal monitor (single collection cycle)
Documentation=file:///opt/thermal_monitor/THERMAL_MONITOR_DESIGN.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# DynamicUser=yes allocates a transient UID/GID for each invocation (in the
# 61184–65519 range), so no `useradd thermal-monitor` setup step is required.
# StateDirectory= files persist across restarts and reboots — systemd
# physically stores them under /var/lib/private/<name>/ (with a symlink at
# /var/lib/<name>) and recursively chowns the tree to the current dynamic
# UID before the service starts, so ownership stays correct even if the
# allocated UID changes between runs.  In practice the same unit name maps
# to a hash-derived UID and reuses it, so chowns are rare.
DynamicUser=yes
WorkingDirectory=/opt/thermal_monitor

# Override these with a drop-in (systemctl edit) if your layout differs.
Environment=CONFIG=/etc/thermal_monitor/thermal_monitor.yaml
Environment=JSON_OUT=/var/lib/thermal_monitor/readings.json
# $THERMAL_MONITOR_STATE_FILE overrides the YAML alerting.state_file setting
# and the code default (/tmp/...).  Pointing it into StateDirectory makes the
# shipped unit work without YAML edits — PrivateTmp=yes below would otherwise
# wipe /tmp state every run and defeat the alert-cooldown logic.
Environment=THERMAL_MONITOR_STATE_FILE=/var/lib/thermal_monitor/state.json
# Don't litter /opt/thermal_monitor with .pyc files — ProtectSystem=strict
# would make that write fail and spam journald with BytecodeError warnings.
Environment=PYTHONDONTWRITEBYTECODE=1

ExecStart=/opt/thermal_monitor/run_monitor.sh \
    -c ${CONFIG} \
    --json ${JSON_OUT} \
    --log-format=systemd

# State / reading-log directory — auto-created and chown'd to the dynamic
# UID.  Visible to both host and service at /var/lib/thermal_monitor (the
# physical path under DynamicUser= is /var/lib/private/thermal_monitor with
# a symlink at the simple location).
StateDirectory=thermal_monitor
StateDirectoryMode=0750

# Logs flow to journald; the sd-daemon <N> prefix on each line supplies the
# priority, so `journalctl -p warning` filters WARN/CRIT correctly.
StandardOutput=journal
StandardError=journal

# Sandboxing.  Most of these are already implied by DynamicUser=yes, but
# we set them explicitly so the policy stays in force if DynamicUser is
# ever removed via a drop-in, and so a security-audit reading this unit
# can see the posture without cross-referencing the systemd manpage.
# ProtectHome=yes is intentionally stricter than DynamicUser's implied
# ProtectHome=read-only — the service has no business reading /home or /root.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes

# Don't restart on failure — the timer will retry on the next cycle.
Restart=no
