# systemd/thermal-monitor.service
# thermal-monitor.service — single collection cycle, triggered by the timer.
#
# Installation:
#   1. Place the project at /opt/thermal_monitor (or override paths below).
#   2. Create a service user:   sudo useradd -r -s /usr/sbin/nologin thermal-monitor
#   3. Edit /etc/thermal_monitor/thermal_monitor.yaml (copy from the example).
#      Point state_file and logging.db_file into /var/lib/thermal_monitor/:
#          alerting:
#            state_file: /var/lib/thermal_monitor/state.json
#          logging:
#            db_file:    /var/lib/thermal_monitor/readings.db
#      (Required: PrivateTmp= gives each run its own /tmp, so cooldown state
#      stored in /tmp would not survive between cycles.)
#   4. Install the units:
#          sudo install -m 644 systemd/thermal-monitor.service /etc/systemd/system/
#          sudo install -m 644 systemd/thermal-monitor.timer   /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now thermal-monitor.timer
#   5. Watch it run:
#          journalctl -u thermal-monitor.service -f
#          journalctl -u thermal-monitor.service -p warning    # only WARN/CRIT
#
# Tune paths / user without editing this file by running:
#      sudo systemctl edit thermal-monitor.service
# which creates a drop-in at /etc/systemd/system/thermal-monitor.service.d/.

[Unit]
Description=Equipment-room thermal monitor (single collection cycle)
Documentation=file:///opt/thermal_monitor/THERMAL_MONITOR_DESIGN.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=thermal-monitor
Group=thermal-monitor
WorkingDirectory=/opt/thermal_monitor

# Override these with a drop-in (systemctl edit) if your layout differs.
Environment=CONFIG=/etc/thermal_monitor/thermal_monitor.yaml
Environment=JSON_OUT=/var/lib/thermal_monitor/readings.json

ExecStart=/opt/thermal_monitor/run_monitor.sh \
    -c ${CONFIG} \
    --json ${JSON_OUT} \
    --log-format=systemd

# State / reading-log directory — auto-created and chown'd to the unit user.
# Available to the process at /var/lib/thermal_monitor.
StateDirectory=thermal_monitor
StateDirectoryMode=0750

# Logs flow to journald; the sd-daemon <N> prefix on each line supplies the
# priority, so `journalctl -p warning` filters WARN/CRIT correctly.
StandardOutput=journal
StandardError=journal

# Sandboxing — keep the blast radius small.  Relax with a drop-in if the
# config or output path falls outside these trees.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes

# Don't restart on failure — the timer will retry on the next cycle.
Restart=no
